# etc/supervisor-britdoc-site-production.ini
[program:britdoc-site-production]
command=/home/ops/.virtualenvs/britdoc-site-production/bin/gunicorn -c /var/www/britdoc-site-production/etc/gunicorn-britdoc-site-production.conf.py --timeout=86400 --graceful-timeout=86400  --log-level=debug  --workers=5 britdoc.wsgi:application
directory=/var/www/britdoc-site-production/
user=ops
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/britdoc-site-production.log
redirect_stderr=true

